# deploy/bankapp.service
[Unit]
Description=Bank API FastAPI Application
After=network.target

[Service]
Type=simple
User=bankapp
Group=bankapp
WorkingDirectory=/opt/bankapp/backend
Environment="PATH=/opt/bankapp/venv/bin"
EnvironmentFile=/opt/bankapp/.env
ExecStart=/opt/bankapp/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=always
RestartSec=5
StandardOutput=append:/var/log/bankapp/app.log
StandardError=append:/var/log/bankapp/error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/bankapp

[Install]
WantedBy=multi-user.target
